# bluefield/misc/DEBIAN/postinst
#!/bin/sh

systemctl enable local_bmc_dev.service
